# living-room-display.service
[Unit]
Description=Living Room Ethereum Display
After=graphical.target
Wants=graphical.target

[Service]
Type=oneshot
User=egk
WorkingDirectory=/home/egk/ethereum-grafana-cluster
Environment=DISPLAY=:0
ExecStart=/home/egk/ethereum-grafana-cluster/start-living-room.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
