# dist/release/living-room-switcher.service
[Unit]
Description=Living Room Display Switcher
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=%U
WorkingDirectory=%h/ethereum-grafana-cluster
Environment=DISPLAY=:0
ExecStart=%h/ethereum-grafana-cluster/start-persistent.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
